# scripts/ansible-dashboard.service
[Unit]
Description=Ansible Playbook Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/ansible/playbooks
ExecStart=/usr/bin/python3 /etc/ansible/playbooks/scripts/log_dashboard.py
Restart=always
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target 